# debian/todoprint.service
[Unit]
Description=GPIO listener for printing todo list
After=network.target

[Service]
Type=simple
ExecStart=/opt/todoprint/gpio_listener.py

[Install]
WantedBy=multi-user.target
